# package/DEBIAN/postrm
#!/bin/bash
set -e

echo "postrm '$1' '$2'"

if [ "$1" = "purge" ]; then
    systemctl disable fleet-telemetry || true
fi
